# puppet/files/gunicorn-supervisord.conf
[program:gunicorn]
command=gunicorn mysite.wsgi:application -c gunicorn.conf.py
directory=/vagrant
user=nobody
autostart=true
autorestart=true
redirect_stderr=true